# wg-ddns-monitor.service
[Unit]
Description=WireGuard DDNS Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wg-ddns-monitor.sh
